# ci-scripts/infra-setup/roles/te-broker/files/te_workers.service
[Unit]
Description=TE Workers
After = network.target

[Service]
ExecStart=/opt/stack/tripleo-ci/scripts/te-broker/start_workers.sh

[Install]
WantedBy = multi-user.target
